# deploy/stavki_bot.service
[Unit]
Description=STAVKI Telegram Bot Service
After=network.target

[Service]
Type=simple
User=serni13678
WorkingDirectory=/home/serni13678/stavki_v2
ExecStart=/usr/bin/python3 /home/serni13678/stavki_v2/scripts/start_bot.py
Restart=always
RestartSec=10
EnvironmentFile=/home/serni13678/stavki_v2/.env
StandardOutput=append:/home/serni13678/stavki_v2/logs/bot_service.log
StandardError=append:/home/serni13678/stavki_v2/logs/bot_error.log

[Install]
WantedBy=multi-user.target
