# deploy/stavki.service
[Unit]
Description=STAVKI Sports Betting Analytics System
After=network.target
Documentation=https://github.com/your-repo/stavki_v2

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/stavki_v2
ExecStart=/usr/bin/python3 -m stavki.interfaces.cli_app scheduler --hours 1

# Environment
EnvironmentFile=/home/%i/stavki_v2/.env
Environment="PYTHONPATH=/home/%i/stavki_v2"

# Auto-restart on crash (30s delay)
Restart=always
RestartSec=30

# Resource limits (prevent runaway)
MemoryMax=2G
CPUQuota=80%

# Logging goes to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stavki

[Install]
WantedBy=multi-user.target
