# deploy/stavki_bot.service
[Unit]
Description=STAVKI v2 Telegram Bot
After=network.target

[Service]
User=root
# Adjust paths as needed for actual user
WorkingDirectory=/home/ubuntu/stavki_v2
ExecStart=/home/ubuntu/stavki_v2/venv/bin/python3 scripts/run_bot.py
Environment="PYTHONPATH=/home/ubuntu/stavki_v2"
Environment="STAVKI_ENV=production"
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
