# deploy/stavki.service
[Unit]
Description=STAVKI Daily Value Bet Pipeline
After=network.target

[Service]
Type=oneshot
User=serni13678
WorkingDirectory=/home/serni13678/stavki_v2
ExecStart=/bin/bash /home/serni13678/stavki_v2/deploy/run_daily.sh
StandardOutput=append:/home/serni13678/stavki_v2/logs/service.log
StandardError=append:/home/serni13678/stavki_v2/logs/service_error.log

[Install]
WantedBy=multi-user.target
